# libvirt/socket_vmnet-shared.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.voipstack.socket_vmnet.shared</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/opt/socket_vmnet/bin/socket_vmnet</string>
        <string>--vmnet-mode=shared</string>
        <string>--vmnet-gateway=192.168.64.1</string>
        <string>--vmnet-dhcp-end=192.168.64.254</string>
        <string>--socket-group=staff</string>
        <string>/opt/homebrew/var/run/socket_vmnet.shared</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/opt/homebrew/var/log/socket_vmnet/shared.stderr</string>
    <key>StandardOutPath</key>
    <string>/opt/homebrew/var/log/socket_vmnet/shared.stdout</string>
</dict>
</plist>
